# launchd/com.hxie.io-worker-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hxie.io-worker-manager</string>

  <key>ProgramArguments</key>
  <array>
    <string>__ROOT__/.venv/bin/python</string>
    <string>__ROOT__/run.py</string>
    <string>--once</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__ROOT__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>18000</integer>

  <key>StandardOutPath</key>
  <string>__ROOT__/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>__ROOT__/logs/launchd.err.log</string>
</dict>
</plist>
